# system/gatewaypi-jack.service
[Unit]
Description=GatewayPi JACK server (owns the USB interface, clocks full-duplex)
After=sound.target systemd-user-sessions.service
Before=gatewaypi-kiosk.service

[Service]
User=@USER@
# Realtime scheduling + locked memory for glitch-free audio.
LimitRTPRIO=95
LimitMEMLOCK=infinity
Environment=JACK_NO_AUDIO_RESERVATION=1
# 48 kHz on the configured interface. @JACKDEV@ (ALSA card, e.g. hw:iTwo)
# and @JACKPERIOD@ (frames/period — 128 on Pi 5, 256 on Pi 4 for headroom)
# are substituted by install.sh based on the detected Pi model.
ExecStart=/usr/bin/jackd -R -P70 -dalsa -d@JACKDEV@ -r48000 -p@JACKPERIOD@ -n3
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
